# salt/common/packages/prometheus/node_exporter/files/node_exporter.service
[Unit]
Description=Node Exporter - Prometheus
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=prometheus
ExecStart={{ install_path }}/node_exporter --web.listen-address=:{{ port }}

[Install]
WantedBy=multi-user.target
